# test/mount/etc/profile.d/mmf.sh
#!/bin/sh
PATH=/bin:/usr/bin:/sbin:/usr/sbin
export GST_DEBUG=2
