# test/mount/usr/lib/systemd/system/bluetooth-hci-down.path
[Unit]
Description=Path activation for the "hciconfig hci0 down"

[Path]
PathExists=/run/.hci_down
