# test/mount/usr/lib/systemd/system/bluetooth-hci-logger.path
[Unit]
Description=Path activation for the "bluetooth-hci-logger"

[Path]
PathExists=/run/.hci_logger_start
