# test/mount/usr/lib/systemd/system/bluetooth-hci-down.service
[Unit]
Description=Bluetooth hci down

[Service]
Type=oneshot
ExecStart=/usr/bin/hciconfig hci0 down
