# test/mount/usr/lib/systemd/system/bluetooth-hci-up.service
[Unit]
Description=Bluetooth hci up

[Service]
Type=oneshot
ExecStart=/usr/bin/hciconfig hci0 up
