# test/mount/usr/lib/systemd/system/bluez-start.path
[Unit]
Description=Path activation for the "bluez"

[Path]
PathExists=/run/.bluez_start
