# test/mount/usr/lib/systemd/system/bluetooth-hci-up.path
[Unit]
Description=Path activation for the "hciconfig hci0 up"

[Path]
PathExists=/run/.hci_up
