# test/mount/usr/lib/systemd/system/bluetooth-hci-logger-stop.path
[Unit]
Description=Path activation for the "bluetooth-hci-logger-stop"

[Path]
PathExists=/run/.hci_logger_stop
