# test/mount/usr/lib/systemd/system/bluetooth-stack-down.path
[Unit]
Description=Path activation for the "bluetooth-stack-down.service"

[Path]
PathExists=/run/.bt_stack_down