# test/mount/usr/lib/systemd/system/device-certificate-manager.target
[Unit]
Description=Device Certificate Manager sockets
DefaultDependencies=true
