# test/mount/usr/lib/systemd/system/display-manager-ready.path
[Unit]
Description=Path activation for display manager ready service
After=tmp.mount display-manager.service

[Path]
PathExists=/run/wayland-0
