# test/mount/usr/lib/systemd/system/dlog_logger.path
[Unit]
Description=Path activation for dlog_logger

[Path]
PathExists=/opt/etc/dlog.conf
