# test/mount/usr/lib/systemd/system/focus-server.path
[Unit]
Description=Path activation for the "focus-server" programme

[Path]
PathExists=/tmp/.pa_ready
