# test/mount/usr/lib/systemd/system/mnt-inform.mount
[Unit]
Description=Reboot Parameter Partition
DefaultDependencies=no
Conflicts=umount.target

[Mount]
What=LABEL=inform
Where=/mnt/inform
Type=ext4
