# test/mount/usr/lib/systemd/system/muse-server.path
[Unit]
Description=muse server

[Path]
PathExists=/run/.wm_ready

