# test/mount/usr/lib/systemd/system/sdbd.socket
[Unit]
Description=sdbd Socket for sdb

[Socket]
ListenUSBFunction=/dev/usb-funcs/sdb/default
Service=sdbd.service
