# test/mount/usr/lib/systemd/system/usb-host-test.socket
[Unit]
Description=USB test functionfs socket

[Socket]
ListenUSBFunction=/run/usb-host-test-ffs
Service=usb-host-ffs-test-daemon.service
