# test/mount/usr/lib/systemd/user/alarm-session-agent.service
[Unit]
Description=Start the alarm agent

[Service]
ExecStart=/usr/bin/alarm_session_agent

[Install]
WantedBy=default.target
