# test/mount/usr/lib/systemd/user/dbus.socket
[Unit]
Description=D-Bus User Message Bus Socket
DefaultDependencies=no

[Socket]
ListenStream=%t/bus

[Install]
Also=dbus.service
