# test/mount/usr/lib/systemd/user/package-recovery-helper.service
[Unit]
Description=Installation Recovery Helper Service

[Service]
ExecStart=/usr/bin/pkg_recovery_helper
