# test/mount/usr/lib/systemd/user/pkgmgr-installer-signal-agent.service
[Unit]
Description=Package Manager Installer Signal Agent on User Session

[Service]
Type=oneshot
ExecStart=/usr/bin/pkgmgr-installer-signal-agent
