# test/mount/usr/lib/systemd/user/tbm-drm-auth-user.path
[Unit]
Description=Wait for tbm-drm-auth display daemon socket

[Path]
PathExists=/run/tbm-drm-auth
