# test/mount/usr/lib/systemd/user/tdm-socket-user.path
[Unit]
Description=Wait for tdm-socket

[Path]
PathExists=/run/tdm-socket
